# agent/crontab
# Gavel poller — detect new Final agenda items every 5 minutes (MOO-41).
*/5 * * * * node /app/scripts/poll-once.mjs
